# scripts/spaces-sync.service
[Unit]
Description=DigitalOcean Spaces FTP Sync Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash /usr/local/bin/spaces_sync_watch.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
